# public/pi-scripts/first-boot-setup.service
[Unit]
Description=Pi Control Center First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/opt/.pi-control-center-installed

[Service]
Type=oneshot
ExecStart=/opt/first-boot-setup.sh
RemainAfterExit=false
TimeoutStartSec=900
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
